# systemd/ezgha-token-refresh.service
[Unit]
Description=Refresh ezgha GitHub App installation token

[Service]
Type=oneshot
ExecStart=@SCRIPTS_DIR@/refresh_gh_app_token.sh
