# systemd/ezgha-token-refresh.service
[Unit]
Description=Refresh ezgha GitHub App installation token

[Service]
Type=oneshot
ExecStart=@REPO_PATH@/scripts/refresh_gh_app_token.sh
