# systemd/ezgha.service.d/10-oomd-omit.conf
# 10-oomd-omit.conf — per-unit systemd-oomd exemption for ezgha.service.
# Part of bead ez-gh-actions-0725, added after adversarial verification
# (sidekick-memarch, 2026-07-10) found that tuning systemd-oomd's
# SwapUsedLimit/DefaultMemoryPressureLimit thresholds tighter (see
# docs/host-ops-sudo-block-0725.md Option A) makes user.slice/
# user-1000.slice's descendants candidates for an IMMEDIATE SIGKILL by
# systemd-oomd — with no cooldown, no grace period, no exclusion list,
# unlike this bead's own user-scope psi-oom-watcher.sh fallback.
#
# VERIFIED ON HOST (2026-07-10, live, read-only): the Colima VM's
# qemu-system-x86_64 process (33.6G RSS) is NOT a separate top-level
# process — it is a CHILD of ezgha.service itself. `cat
# /proc/<qemu-pid>/cgroup` returns
# `0::/user.slice/user-1000.slice/user@1000.service/app.slice/ezgha.service`,
# and `systemctl --user status ezgha.service` confirms limactl/qemu as
# descendants under that exact cgroup, with the service's own reported
# Memory: 33.6G (peak 36.3G) — i.e. the entire Colima VM's memory footprint
# is accounted INSIDE ezgha.service's cgroup, not a separate slice. That
# means tightening thresholds on any ANCESTOR slice (user.slice,
# user-1000.slice) without an explicit per-unit exemption on ezgha.service
# itself would make the VM a live kill candidate the moment aggregate
# pressure/swap in that subtree crosses the tightened threshold — exactly
# the failure mode this whole bead exists to prevent, but WORSE (SIGKILL,
# instant, no cooldown) than the incident it's meant to guard against.
#
# ManagedOOMPreference=omit tells systemd-oomd to never assess this unit's
# cgroup (or its descendants — i.e. the qemu/colima subtree) as a kill
# candidate, and to exclude its resource usage from ancestor rollups used
# for candidate selection. Per systemd.resource-control(5): the omit/avoid
# extended attribute is only respected when "the unit's cgroup owner and
# the owner of the monitored ancestor cgroup are the same" — both
# ezgha.service and user-1000.slice are owned by this UID (confirmed via
# the delegated cgroup.controllers check performed earlier in this bead),
# so this exemption is expected to be honored without needing root.
#
# THIS FILE REQUIRES NO SUDO: ezgha.service is a systemd --user unit, so
# this override lives under ~/.config/systemd/user/ezgha.service.d/, which
# the invoking user owns outright. It is NOT auto-installed by install.sh
# (single-writer rule: no live systemctl actions were run as part of this
# commit) — see docs/host-ops-sudo-block-0725.md for the manual install +
# required restart step (the omit xattr is only applied to a unit's live
# cgroup when the unit next starts, so this requires an ezgha.service
# restart to take effect after installing this drop-in — coordinate with
# whoever the ezgha deploy-owner is per repo CLAUDE.md single-writer rule,
# do not restart opportunistically).
[Service]
ManagedOOMPreference=omit
