# systemd/ao-daemon-memory-cap.service.d/memory.conf
# ao-daemon-memory-cap.service.d/memory.conf — tracked drop-in that, when
# copied into ~/.config/systemd/user/, caps the ao-daemon (Agent
# Orchestrator) inside agents.slice with MemoryHigh=20G. Part of bead
# ez-gh-actions-0725 (Round 3 lane G, AO/agent auto-migration).
#
# INSTALL (single-writer only — never dispatch a concurrent install while
# the parent sidekick is also touching this directory):
#   # 1) ensure ao-daemon.service is at the canonical path
#   [ -f ~/.config/systemd/user/ao-daemon.service ] || \
#       ln -sf <repo>/systemd/ao-daemon.service ~/.config/systemd/user/ao-daemon.service
#   # 2) create the drop-in folder + this conf file
#   mkdir -p ~/.config/systemd/user/ao-daemon.service.d
#   cp systemd/ao-daemon-memory-cap.service.d/memory.conf \
#      ~/.config/systemd/user/ao-daemon.service.d/memory.conf
#   # 3) reload + restart so the new slice assignment takes effect
#   systemctl --user daemon-reload
#   systemctl --user restart ao-daemon.service
#
# After install, `cat /sys/fs/cgroup$(cat /proc/$(pgrep -f ao-daemon |
# head -1)/cgroup | sed 's/^0:://')/memory.high` should report a finite
# value (not "max"), and Gate 8 (2) of docs/verify-exit-criteria.sh
# will treat the daemon as enrolled.
#
# PLACEHOLDER SERVICE NAME: this drop-in's [Install] section is
# intentionally absent (drop-ins are unit-agnostic — they merge into
# whichever unit is loaded at the matching name). The
# `Slice=` directive enrolls the daemon's main PID into agents.slice
# regardless of the service's name, so the conf file works whether the
# service is named `ao-daemon.service` (canonical), `ao.service`, or
# anything downstream renames it to. If your install uses a non-
# canonical name (e.g. `ao-orchestrator.service`), copy this file to
# THAT service's drop-in folder instead.
#
# WHY Slice=agents.slice (rather than allocating a sibling slice): the
# round-3 panel decision is that EVERY interactive agent CLI process —
# claude, codex, gemini, cursor, aider, cody, AND the AO daemon that
# spawns them — must share one bounded budget so their aggregate RAM
# demand cannot exceed the slice's 20G ceiling. Putting the daemon in
# its own slice would defeat the aggregate-cap invariant and silently
# re-introduce the 2026-07-10 unbounded-daemon failure mode.
#
# BLAST RADIUS: MemoryHigh=20G is a SOFT ceiling (throttles and reclaims
# page cache; does NOT OOM-kill). The 2026-07-12 measured peak of the
# AO daemon workload is 2.4-2.7 GiB live RSS × ~382 tasks; 20G leaves
# ~17G headroom. Slice=agents.slice inherits the parent slice's
# MemoryHigh if this drop-in is ever installed without the agents.slice
# unit file in place — but systemd will silently create an ad-hoc
# unbounded slice in that case, so the install instructions above
# explicitly gate on the slice unit file.
#
# Atomic update: this file is one-shot-tracked; do NOT regenerate it
# automatically. A future enhancement (bead candidates welcome) could
# template MemoryHigh off `total_host_RAM * 0.25` rather than the
# hard-coded 20G.
[Service]
MemoryHigh=20G
MemoryAccounting=yes
Slice=agents.slice
